# mint-security-update.service
[Unit]
Description=Mint security updates

[Service]
ExecStart=mintupdate-cli -sy upgrade
IOWeight=20
